# infra/systemd/openclaw-bridge.service
[Unit]
Description=OpenClaw WhatsApp Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bridge
EnvironmentFile=/home/ubuntu/bridge/.env
ExecStart=/usr/bin/node /home/ubuntu/bridge/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
